# docker/s6/.s6-svscan/finish
#!/bin/sh

# Cleanup SOCAT services and s6 event folder
rm -rf $(find /app/asouldocs/docker/s6/ -name 'event')
rm -rf /app/asouldocs/docker/s6/SOCAT_*
