# docker/s6/peach/run
#!/bin/sh

if test -f ./setup; then
    source ./setup
fi

export USER=asouldocs
exec gosu $USER /app/asouldocs/asouldocs web
