# docker/s6/peach/run
#!/bin/sh

if test -f ./setup; then
    source ./setup
fi

export USER=peach
exec gosu $USER /app/peach/peach web
